# crontab

0,5,10,15,20,25,30,35,40,45,50,55 * * * * /root/noaaweather.pl > /dev/null 2>&1

0 * * * * nice /root/wapleslapse.sh >/dev/null 2>&1
